# systemd/buzz-relay.service
[Unit]
Description=Buzz Relay Server
Documentation=https://github.com/block/buzz
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=marvin
WorkingDirectory=/home/marvin/.openclaw/workspace/repos/buzz
EnvironmentFile=/home/marvin/.openclaw/workspace/repos/buzz/.env
ExecStart=/home/marvin/.local/bin/buzz-relay
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=buzz-relay

# Give the relay time to connect to Postgres/Redis/MinIO on startup
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
